# .config/roles/systemd_container/files/Debian-11/etc/systemd/system/serial-getty@.service.d/noclear.conf
# 2016-2021 Jakob Meng, <jakobmeng@web.de>
# Disable clearing of boot messages
# Reference: https://wiki.archlinux.org/index.php/disable_clearing_of_boot_messages

[Service]
TTYVTDisallocate=no
