# files/Debian-10/lib/systemd/system/prepare-container.service
# (2018-2019) WEARENOTALONE

[Unit]
Description=prepare-container
#DefaultDependencies=no
After=local-fs.target
ConditionFileIsExecutable=/usr/bin/prepare-container

[Service]
Type=oneshot
PassEnvironment=DISPLAY
ExecStart=/usr/bin/prepare-container
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=basic.target
