# debian-testing-kde-dev/prepare-container.service.d/override.conf
# (2018) WEARENOTALONE

[Service]
PassEnvironment=XDG_CURRENT_DESKTOP
