# .config/roles/systemd_container/files/Debian-10/etc/systemd/system/prepare-container.service.d/override.conf.example
[Service]
PassEnvironment=EXAMPLE_ENVIRONMENT_VARIABLE
